# services/sys-manager-cleanup.timer
[Unit]
Description=Run weekly cleanup of old system manager backups

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
